# configuration/roles/backup/templates/backup.cron.j2
# Backup cron jobs - runs daily at 3 AM
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin

# Panel backup (PostgreSQL)
0 3 * * * root docker run --rm \
--env-file /opt/remnawave-backups/scripts/backup.env \
-v /var/run/docker.sock:/var/run/docker.sock \
-v /opt/remnawave-backups/data:/opt/remnawave-backups/data \
remnawave-backup:latest backup --panel >> /var/log/remnawave-backup.log 2>&1